# deploy/systemd/raceframe-maintenance.timer
[Unit]
Description=Run RaceFrame maintenance every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
RandomizedDelaySec=20s

[Install]
WantedBy=timers.target
